# board/pumpkin-mbm2/rootfs_overlay/etc/network/if-up.d/arp_sparcam.sh
arp -i eth0 -s 200.1.1.101 3c:18:a0:07:eb:32

